# deploy/systemd/btc_trade_bot.service.example
[Unit]
Description=BTC Demo Trade Bot (Headless)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/btc_trade_bot/BTC_trade_bot
EnvironmentFile=/opt/btc_trade_bot/.env
ExecStart=/opt/btc_trade_bot/venv/bin/python /opt/btc_trade_bot/BTC_trade_bot/live_trading_bot.py --headless --risk low
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30
StandardOutput=append:/var/log/btc_trade_bot.log
StandardError=append:/var/log/btc_trade_bot.log

[Install]
WantedBy=multi-user.target
